# launchd/com.openclaw.xchat.orchestrator.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.openclaw.xchat.orchestrator</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/node</string>
      <!-- Update this path to your clone location -->
      <string>/ABSOLUTE/PATH/TO/clawdbot-xchat-orchestrator/src/orchestrator.mjs</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>${HOME}</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>${HOME}</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>${HOME}/.openclaw/xchat/orchestrator.log</string>

    <key>StandardErrorPath</key>
    <string>${HOME}/.openclaw/xchat/orchestrator.err.log</string>
  </dict>
</plist>
